# supervisor.ini
[program:youtube]
directory=/deploy
command=python3 main.py
autostart=true
autorestart=unexpected
startretries=3
stopasgroup=true
killasgroup = true
stdout_logfile=/tmp/youtube.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1
redirect_stderr=true